# infra/var_gold.service
[Unit]
Description=var_gold PAXG/XAUT monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/var_gold
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/var_gold/var_spread_monitor.py --log-level INFO
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
